# contrib/init/vibofchaind.conf
description "VibofChain Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env VIBOFCHAIND_BIN="/usr/bin/vibofchaind"
env VIBOFCHAIND_USER="vibofchain"
env VIBOFCHAIND_GROUP="vibofchain"
env VIBOFCHAIND_PIDDIR="/var/run/vibofchaind"
# upstart can't handle variables constructed with other variables
env VIBOFCHAIND_PIDFILE="/var/run/vibofchaind/vibofchaind.pid"
env VIBOFCHAIND_CONFIGFILE="/etc/vibofchain/vibofchain.conf"
env VIBOFCHAIND_DATADIR="/var/lib/vibofchaind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # vibofchaind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$VIBOFCHAIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run vibofchaind."
        echo "The setting must appear in $VIBOFCHAIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"VibofChain Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$VIBOFCHAIND_PIDDIR"
    chmod 0755 "$VIBOFCHAIND_PIDDIR"
    chown $VIBOFCHAIND_USER:$VIBOFCHAIND_GROUP "$VIBOFCHAIND_PIDDIR"
    chown $VIBOFCHAIND_USER:$VIBOFCHAIND_GROUP "$VIBOFCHAIND_CONFIGFILE"
    chmod 0660 "$VIBOFCHAIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$VIBOFCHAIND_PIDFILE" \
    --chuid $VIBOFCHAIND_USER:$VIBOFCHAIND_GROUP \
    --exec "$VIBOFCHAIND_BIN" \
    -- \
    -pid="$VIBOFCHAIND_PIDFILE" \
    -conf="$VIBOFCHAIND_CONFIGFILE" \
    -datadir="$VIBOFCHAIND_DATADIR" \
    -disablewallet \
    -daemon

